# mule4/EE-license/resources/am-2.5.11.2/am/config/am-filebeat-watcher.template.path
[Path]
PathChanged=$muleHome$/am/config/filebeat.yml
PathChanged=$muleHome$/am/config/certs/mule-agent-chain.crt
PathChanged=$muleHome$/am/config/certs/mule-agent.key
PathChanged=$muleHome$/am/config/certs/registration-agent-ca.crt

[Install]
WantedBy=multi-user.target
