# mule4/EE-license/resources/am-2.5.11.2/am/config/am-filebeat-watcher.template.service
[Unit]
Description=Anypoint Monitoring Filebeat Watcher Service
Documentation=https://github.com/mulesoft/dias-agent-installer
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=$systemctlExec$ restart am-filebeat.service

[Install]
WantedBy=multi-user.target
